# config/inittab
::sysinit:/etc/init.d/rcS
::restart:/sbin/init
::shutdown:/bin/umount -a -r
ttyS0::askfirst:-/bin/sh

